# catty-app.service
[Unit]
Description=Catty Reminders App
After=network.target

[Service]
Type=simple
User=deemeed
Group=deemeed
WorkingDirectory=/home/deemeed/catty-reminders-app
EnvironmentFile=/home/deemeed/catty-reminders-app/.env
ExecStart=/home/deemeed/catty-reminders-app/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8181
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
