# catty.service
[Unit]
Description=Catty FastAPI Reminders App
After=network.target

[Service]

User=yui1337
Group=www-data

WorkingDirectory=/home/yui1337/my_project

EnvironmentFile=/home/yui1337/my_project/.env

ExecStart=/home/yui1337/my_project/venv/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8181

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
