# devops-webhook.service
[Unit]
Description=DevOps Webhook Handler (port 8080)
After=network.target frpc.service

[Service]
User=kirill
WorkingDirectory=/home/kirill/desktop/devops
ExecStart=/usr/bin/python3 /home/kirill/desktop/devops/webhook-server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
