# webhook-handler.service
[Unit]
Description=Webhook Handler for GitHub
After=network.target

[Service]
Type=simple
User=vano
WorkingDirectory=/home/vano/catty-reminders-app
ExecStart=/home/vano/catty-reminders-app/venv/bin/python /home/vano/catty-reminders-app/webhook_handler.py
Restart=always

[Install]
WantedBy=multi-user.target
