# webhook-handler.service
[Unit]
Description=Webhook Handler
After=network.target

[Service]
User=vano
WorkingDirectory=/home/vano/catty-reminders-app
ExecStart=/usr/bin/python3 /home/vano/catty-reminders-app/webhook_handler.py
Restart=always

[Install]
WantedBy=multi-user.target
