# webhook-handler.service
[Unit]
Description=Webhook Handler for DevOps Lab
After=network.target

[Service]
Type=simple
User=andrey
WorkingDirectory=/home/andrey/Desktop/DevOps_lab1/catty-reminders-app
ExecStart=/usr/bin/python3 /home/andrey/Desktop/DevOps_lab1/catty-reminders-app/webhook_handler.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
