# webhook-handler.service
[Unit]
Description=Webhook Handler
After=network.target

[Service]
User=user
WorkingDirectory=/home/user/webhook-handler
ExecStart=/home/user/webhook-handler/venv/bin/python logic.py
Restart=always

[Install]
WantedBy=multi-user.target
