# systemd/catty-app.service
[Unit]
Description=Catty Reminders App
After=network.target

[Service]
User=kali
WorkingDirectory=/home/kali/catty-reminders-app
EnvironmentFile=-/home/kali/catty-reminders-app/.env.deploy
ExecStart=/home/kali/catty-reminders-app/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8181
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
