# systemd/webhook-handler.service
[Unit]
Description=Webhook Handler
After=network.target

[Service]
User=password_123
WorkingDirectory=/home/password_123/webhook-handler
ExecStart=/home/password_123/webhook-handler/venv/bin/python script.py
Restart=always

[Install]
WantedBy=multi-user.target
