# deploy/webhook.service
[Unit]
Description=Github Webhook Listener
After=network.target

[Service]
User=vanya
Group=vanya
WorkingDirectory=/home/vanya/webhook-server
# Предполагаю, что ты используешь системный питон или venv в этой папке
ExecStart=/usr/bin/python3 app.py
Restart=always

[Install]
WantedBy=multi-user.target
